# opt/virtualbox/.config/dotfiles/etc/profile.d/virtalbox.sh
#!/bin/sh
export VBOX_USER_HOME="${XDG_CONFIG_DIR}/virtualbox"
