# opt/pip/.config/dotfiles/etc/profile.d/pip.sh
#!/bin/sh
export PIP_CONFIG_FILE=${XDG_CONFIG_HOME}/pip/conf
export PIP_LOG_FILE=${XDG_CACHE_HOME}/pip/log
