# opt/screen/.config/dotfiles/etc/profile.d/screen.sh
#!/bin/sh
export SCREENRC="${XDG_CONFIG_HOME}/screen/rc"
